# services/feh.service
[Unit]
Description=BackgroundChanger

[Service]
ExecStart=/home/grahnen/.dotfiles/scripts/bg.sh

[Install]
WantedBy=default.target
